# core.env/Library/LaunchAgents/env.plist
<?xml version="1.0" encoding="UTF-8" ?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>my.environment</string>
  <key>ProgramArguments</key>
  <array>
    <string>sh</string>
    <string>-c</string>
    <string>
      launchctl setenv XDG_CONFIG_HOME /Users/tjgorton/.config
      launchctl setenv XDG_DATA_HOME /Users/tjgorton/.local/share
      launchctl setenv XDG_STATE_HOME /Users/tjgorton/.local/state
      launchctl setenv ZDOTDIR /Users/tjgorton/.config/zsh
    </string>
  </array>
  <key>RunAtLoad</key>
  <true />
</dict>
</plist>
